# php/laravel/mysql57/cron/resource/cron/crontabs/root
# e.g.
* * * * * docker exec php73mysql57laravel_app_1 php artisan schedule:run
